# pipeline/server_utils/backup_apps.timer
[Unit]
Description=Run weekly app data backup

[Timer]
OnCalendar=daily 03:00:00
Persistent=true

[Install]
WantedBy=timers.target

