# pipeline/server_utils/backup_apps.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/sbin/btrbk -c /etc/btrbk/apps.conf run
