# pipeline/apiserv/xray-apiserv.service
[Unit]
Description=The API server for the XRay project

[Service]
ExecStart=/usr/local/bin/apiserv
User=xray
Group=xray
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
